# systemd/web-terminal.service
[Unit]
Description=Web Terminal
After=network.target

[Service]
Type=simple
User=__WT_USER__
Group=__WT_GROUP__
WorkingDirectory=__WT_PROJECT_DIR__
Environment=HOME=__WT_HOME__
Environment=SHELL=__WT_SHELL__
Environment=PATH=__WT_PATH__
Environment=PORT=3456
Environment=WEB_TERMINAL_DATA_DIR=__WT_PROJECT_DIR__/.web-terminal
ExecStartPre=/usr/bin/mkdir -p __WT_PROJECT_DIR__/.web-terminal/logs
ExecStart=__WT_PROJECT_DIR__/scripts/start-web-terminal.sh
Restart=always
RestartSec=3
StandardOutput=append:__WT_PROJECT_DIR__/.web-terminal/logs/service.log
StandardError=append:__WT_PROJECT_DIR__/.web-terminal/logs/service.log

[Install]
WantedBy=multi-user.target
